# supervisord.conf
[supervisord]
nodaemon = true

[program:livemeasurement]
command =/usr/local/bin/run.sh
startsecs = 2
autorestart = true
startretries = 5
autostart=true
